# templates/crontab
0 11 * * 1-5 python3 ~/scripts/d0check.py
